# base/etc/cron.daily/salt
#!/bin/sh
# update minion

if /sbin/rc-service -qq -e salt-minion
then
	/usr/bin/salt-call state.highstate
fi
