# etc/udev/rules.d/90-display.rules
ACTION=="change", SUBSYSTEM=="drm", RUN+="display.sh &"
